# clawconductor.service
[Unit]
Description=ClawConductor Escalation Proxy
After=network.target litellm.service

[Service]
Type=simple
User=matt
WorkingDirectory=/home/matt/clawconductor
ExecStart=/home/matt/clawconductor/.venv/bin/uvicorn clawconductor.proxy:app --port 8765 --host 127.0.0.1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
